# deploy/systemd/audiobook-organizer.service
# file: deploy/systemd/audiobook-organizer.service
# version: 4.0.0
# guid: f1e2d3c4-b5a6-9870-8765-4321fedcba98
#
# systemd service unit for audiobook-organizer on Linux
#
# Installation:
#   1. Create service user:
#      sudo useradd -r -s /usr/sbin/nologin -d /var/lib/audiobook-organizer audiobook
#   2. Create and set up directories:
#      sudo mkdir -p /var/lib/audiobook-organizer
#      sudo mkdir -p /var/log/audiobook-organizer
#      sudo chown audiobook:audiobook /var/lib/audiobook-organizer
#      sudo chown audiobook:audiobook /var/log/audiobook-organizer
#   3. Install binary:
#      sudo cp audiobook-organizer /usr/local/bin/
#      sudo chmod 0755 /usr/local/bin/audiobook-organizer
#   4. Install service file:
#      sudo cp deploy/systemd/audiobook-organizer.service /etc/systemd/system/
#      sudo systemctl daemon-reload
#   5. Enable and start:
#      sudo systemctl enable --now audiobook-organizer
#
# File access:
#   The audiobook service user needs read access to your library. Configure with:
#   Option A (group membership):
#     sudo usermod -aG media audiobook
#   Option B (ACLs, for more fine-grained control):
#     sudo setfacl -R -m u:audiobook:rX /path/to/audiobooks
#
# Management:
#   sudo systemctl start audiobook-organizer
#   sudo systemctl stop audiobook-organizer
#   sudo systemctl restart audiobook-organizer
#   sudo systemctl status audiobook-organizer
#   journalctl -u audiobook-organizer -f

[Unit]
Description=Audiobook Organizer
Documentation=https://github.com/jdfalk/audiobook-organizer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=audiobook
Group=audiobook

# Service execution — adjust AUDIOBOOK_ROOT_DIR for your library path
ExecStart=/usr/local/bin/audiobook-organizer serve --host 0.0.0.0 --port 8484

# Working directory and environment configuration
WorkingDirectory=/var/lib/audiobook-organizer
Environment="HOME=/var/lib/audiobook-organizer"
Environment="AUDIOBOOK_ROOT_DIR=/var/lib/audiobooks"
Environment="DATABASE_PATH=/var/lib/audiobook-organizer/audiobooks.pebble"

# Restart policy: restart on failure with 5-second delay
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Security hardening — kept minimal to allow reading from various mount points
# (NFS, external drives, etc). The service needs:
#   - Read access to audiobook source directories (configured at AUDIOBOOK_ROOT_DIR)
#   - Read/write to /var/lib/audiobook-organizer (database and metadata)
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
PrivateTmp=yes

# Logging to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audiobook-organizer

[Install]
WantedBy=multi-user.target
