# deploy/launchd/com.jdfalk.audiobook-organizer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- file: deploy/launchd/com.jdfalk.audiobook-organizer.plist -->
<!-- version: 2.0.0 -->
<!-- guid: 1f2e3d4c-5b6a-4789-90ab-cdef12345678 -->
<!--
  macOS launchd plist for audiobook-organizer

  Installation:
    1. Edit this file and replace /Users/USERNAME with your actual username
    2. Create log directory: mkdir -p ~/Library/Logs
    3. Copy to: cp deploy/launchd/com.jdfalk.audiobook-organizer.plist ~/Library/LaunchAgents/
    4. Load: launchctl load ~/Library/LaunchAgents/com.jdfalk.audiobook-organizer.plist
    5. Verify: launchctl list | grep audiobook-organizer

  Management:
    # Start the service
    launchctl start com.jdfalk.audiobook-organizer

    # Stop the service
    launchctl stop com.jdfalk.audiobook-organizer

    # Unload (disable startup)
    launchctl unload ~/Library/LaunchAgents/com.jdfalk.audiobook-organizer.plist

    # View logs
    tail -f ~/Library/Logs/audiobook-organizer.log
    tail -f ~/Library/Logs/audiobook-organizer.err
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jdfalk.audiobook-organizer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/audiobook-organizer</string>
        <string>serve</string>
    </array>

    <!-- Environment variables for runtime configuration -->
    <key>EnvironmentVariables</key>
    <dict>
        <!-- Root directory where audiobook files are located (EDIT THIS) -->
        <key>AUDIOBOOK_ROOT_DIR</key>
        <string>/Users/USERNAME/Audiobooks</string>
        <!-- Path to application database (EDIT THIS) -->
        <key>DATABASE_PATH</key>
        <string>/Users/USERNAME/.audiobook-organizer/db</string>
        <!-- Web server port -->
        <key>PORT</key>
        <string>8484</string>
    </dict>

    <!-- Working directory for relative paths -->
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME/.audiobook-organizer</string>

    <!-- Start at login and restart on failure -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Standard output and error logging -->
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/audiobook-organizer.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/audiobook-organizer.err</string>

    <!-- Prevent multiple instances, restrict file creation permissions -->
    <key>Umask</key>
    <integer>0077</integer>
</dict>
</plist>
